# CLI Pulse Bar/CLI Pulse Bar/HelperAgent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  LaunchAgent plist for the embedded CLI Pulse helper daemon.

  Phase 4D P1.4 (Codex review): this plist is shipped INSIDE the
  signed app bundle at `Contents/Library/LaunchAgents/yyh.CLI-Pulse.
  helper.plist` and MUST NOT be mutated at runtime — any post-
  distribution write into a signed bundle invalidates the code
  signature and causes SMAppService.register to fail.

  Three resolutions baked in below to make the plist stable
  across all install paths without runtime substitution:

    1. ProgramArguments uses `BundleProgram` keyed entry instead
       of `Program` — launchd resolves this RELATIVE TO the
       containing .app bundle automatically. We don't need to
       know the user's install path.

    2. Logging goes via `Logger(subsystem: "yyh.CLI-Pulse.helper.agent",
       ...)` os_log — visible via `log show --predicate
       "process == 'cli_pulse_helper'"`. The plist had
       `StandardOutPath` / `StandardErrorPath` keys with `~/Library/
       Logs/CLI Pulse/...` paths previously, but launchd does NOT
       expand `~` in those keys (Phase 4E e2e 2026-05-07 finding):
       it tried to literally create `//~/Library/Logs/...` in `/`
       (read-only) → exits 78 (EX_CONFIG) before main runs. Keys
       removed; helper logs through the unified os_log subsystem
       which doesn't need a writable file path.

    3. Supabase URL + anon key are NOT in the plist anymore.
       The Swift helper reads them from the app's Info.plist via
       a small lookup at startup (the app's bundle path is
       reachable through `$XPC_BUNDLE_PATH` launchd sets, OR via
       the helper-config file the app writes once on first
       launch). Either way, no secrets cross the plist boundary
       and the LaunchAgent file works for every Mac/install
       location identically.

  The agent runs OUTSIDE the macOS app sandbox (LaunchAgents under
  user-level launchd are not sandboxed by default), which is the
  whole point of this approach — the helper needs to:

    * spawn `claude` PTY children (sandboxed apps can't do that
      with full functionality preserved)
    * read `~/.claude/settings.json` (outside group container)
    * walk the user's git project directories for the activity
      collector

  Communication with the sandboxed app stays via the existing app-
  group container UDS socket at
  `~/Library/Group Containers/group.yyh.CLI-Pulse/clipulse-helper.sock`,
  which is a regular filesystem path both processes can access.

  Bundle identifier scheme:
    * The agent's Label MUST match the plist filename minus `.plist`,
      so `yyh.CLI-Pulse.helper.agent.plist` → label `yyh.CLI-Pulse.helper.agent`.
      Renamed from `yyh.CLI-Pulse.helper` to disambiguate from the
      MAS LoginItem (CLIPulseHelper.app, bundle ID `yyh.CLI-Pulse.helper`)
      which previously claimed the same launchd label slot.
    * SMAppService.agent(plistName:) resolves the plist relative to
      the calling app's Contents/Library/LaunchAgents/ — no absolute
      path needed at the SMAppService API surface.

  Lifecycle:
    * `RunAtLoad`: false — we don't auto-start on machine boot;
      registration happens lazily when the user opens the macOS app
      and the Sessions feature actually needs the helper.
    * `KeepAlive`: true — if the helper crashes, launchd restarts it
      so the UDS socket comes back and the app's poll loop reconnects.
    * `ProcessType`: Background — lower CPU/memory priority; we're a
      polling daemon, not a foreground UI process.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>yyh.CLI-Pulse.helper.agent</string>

    <!--
      Phase 4D P1.4: BundleProgram is the launchd-blessed way to
      reference an executable INSIDE an .app bundle without
      hard-coding an install-path-dependent absolute path. The
      `BundleProgram` value is interpreted relative to the
      registering app's bundle root (i.e. the `.app/` dir).
      `ProgramArguments` provides argv[0] and the trailing
      arguments — `cli_pulse_helper` here is just the argv[0]
      label, not a path. launchd uses BundleProgram for the
      executable resolution.
    -->
    <key>BundleProgram</key>
    <string>Contents/Helpers/cli_pulse_helper</string>
    <key>ProgramArguments</key>
    <array>
        <string>cli_pulse_helper</string>
        <string>daemon</string>
        <string>--interval</string>
        <string>120</string>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <true/>

    <key>ProcessType</key>
    <string>Background</string>

    <!--
      Throttle restart loop on repeated crash. Default 10s is fine
      for a healthy helper but pathological crash loops would log-
      spam without this. 30s gives the user enough time to see the
      Sessions tab go red before the next attempt.
    -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <!--
      No `StandardOutPath` / `StandardErrorPath` — see the comment
      block at top. launchd doesn't expand `~` in those keys, and
      hardcoding an absolute path would break across users / Macs.
      Helper logging is through `Logger(subsystem: "yyh.CLI-Pulse.
      helper", ...)` (os_log), readable via `log show --predicate
      "process == 'cli_pulse_helper'"`.
    -->
</dict>
</plist>
